# trunk/support/rc.local
#!/bin/sh

clear
PATH=$PATH:/usr/local/bin/

#bring the interfaces up
for int in `ifconfig -l`
do
 ifconfig $int up
done

#


cat /etc/motd

/bin/sh 

echo Shutting down

halt -p
