# trunk/support/rc.local
#!/bin/sh

clear

#bring the interfaces up
for int in `ifconfig -l`
do
 ifconfig $int up
done

#echo ""
#echo ""
#echo ""
#echo ""
#echo ""
#echo ""
#echo ""
#echo ""
#echo ""
echo ""
echo "Starting FreeIPS in verbose mode.. use CTRL-C to quit"
echo "then type 'FreeIPS' to start the IDS manually"
echo ""
#
#
cd /usr/local/share/FreeIPS/ && FreeIPS -c config/config.xml 


